# conf/letsencrypt/crontab
30 2 * * 1 /opt/letsencrypt/%(renew_script)s >> /var/log/le-renewal.log
