# conf/letsencrypt/crontab
30 2 * * 1 /opt/letsencrypt/%(renew_script) >> /var/log/le-renewal.log